# deployment/btc-bot.service
[Unit]
Description=BTC Pump/Dump Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=btcbot
Group=btcbot
WorkingDirectory=/opt/btc-bot
Environment="PATH=/opt/btc-bot/venv/bin"
ExecStart=/opt/btc-bot/venv/bin/python /opt/btc-bot/main.py

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/btc-bot

# Resource limits
MemoryMax=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc-bot

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
